# SoundFx/server_sound.conf
[program:motion_skull]
command=/usr/bin/python3 /home/pi/ServerSound.py
directory=/home/pi/
user=pi
autostart=true
autorestart=true